# chipcam.service
[Unit]
Description=camera

[Service]
Type=simple
ExecStart=/home/chip/record.sh

[Install]
WantedBy=multi-user.target
